# .config/systemd/user/take_a_break.service
[Unit]
Description=Take a break for 5 minutes

[Service]
Type=simple
ExecStart=%h/script/bash/take_a_break

[Install]
WantedBy=default.target
